# eink_layout_manager/rootfs/etc/services.d/app/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting eInk Layout Manager frontend & backend..."

export INGRESS_PORT=$(bashio::addon.ingress_port)

# Ensure /data is writable by the addon user (it's a mounted volume)
chown -R addon:addon /data

# Run the python backend as the addon user
exec s6-setuidgid addon python3 /app/main.py
