# eink_layout_manager/rootfs/etc/services.d/mqtt/run
#!/usr/bin/with-contenv bashio

bashio::log.info "Starting MQTT Manager..."

if bashio::services.available "mqtt"; then
    export MQTT_HOST=$(bashio::services mqtt "host")
    export MQTT_PORT=$(bashio::services mqtt "port")
    export MQTT_USER=$(bashio::services mqtt "username")
    export MQTT_PASSWORD=$(bashio::services mqtt "password")
else
    bashio::log.warning "MQTT service not available, MQTT Manager might not connect."
fi

export PYTHONPATH=/
exec s6-setuidgid addon python3 -m backend.mqtt_manager
